# rootfs/etc/services.d/influxdb/run
#!/usr/bin/with-contenv bash
echo "[INFO] Starting InfluxDB v3"
exec /usr/bin/influxdb3 --config /app/influx.toml