# rootfs/etc/services.d/node/run
#!/usr/bin/with-contenv bash
echo "[INFO] Starting Node.js WUHAa server"
exec node /app/index.js
